# init_scripts/etc/service/flannel/run
#!/bin/sh
exec 2>&1
exec flanneld --etcd-endpoints=http://etcd:2379

